# GATEWAY/zigbee-agent/greenhouse_agent.service
[Unit]
Description = Plant green house agent Service
After = rc.local.service

[Service]
Type = simple
User = root
Group = root
ExecStart = /usr/bin/python /usr/local/src/zigbee-agent/greenhouse_agent.py

[Install]
WantedBy = multi-user.target
